# deploy_tools/gunicorn.template.socket
[Unit]
Description=gunicorn SITENAME socket

[Socket]
ListenStream=/run/gunicorn.SITENAME/socket

[Install]
WantedBy=sockets.target
